# deploy/agent.service
[Unit]
Description=AI Agent Docker Runtime
After=network-online.target docker.service tailscaled.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=agent
Group=agent
WorkingDirectory=/opt/agent
ExecStart=/usr/bin/docker compose --env-file deploy/docker/.env up -d --build
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=0
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
